# portfolio_bot.service
[Unit]
Description=Portfolio Bot Service
After=network.target

[Service]
User=deploy
Group=deploy
WorkingDirectory=/opt/bots/portfolio_bot
ExecStart=/opt/bots/portfolio_bot/venv/bin/python3 -m app.main
Restart=always
RestartSec=10
EnvironmentFile=/opt/bots/portfolio_bot/.env

[Install]
WantedBy=multi-user.target
